# deployments/dpkg/ottersacle/DEBIAN/prerm
#!/bin/bash

set -e

# Stop and disable service before removal
systemctl stop ottersacle.service 2>/dev/null || true
systemctl disable ottersacle.service 2>/dev/null || true

exit 0
